# dyndhcpd@.service
[Unit]
Description=IPv4 DHCP server for interface %i
After=network.target

[Service]
ExecStart=/usr/bin/dyndhcpd -i%i

[Install]
WantedBy=multi-user.target
